# crontab
SHELL=/bin/bash
BASH_ENV=/container.env
0 11 * * * root cd /app && ./run.sh > /proc/1/fd/1 2>/proc/1/fd/2
# must be ended with a new line "LF" (Unix) and not "CRLF" (Windows)
